# custom-modules/mailhog/files/mailhog.service
[Unit]
Description=Start Mailhog Web Server
After=network.target

[Service]
ExecStart=/usr/local/bin/mailhog

[Install]
WantedBy=multi-user.target
